# deploy/fly/supervisord.conf
[program:python]
command=/opt/venv/bin/uvicorn app:app --host 127.0.0.1 --port 5001 --workers %(ENV_UVICORN_WORKERS)s
directory=/app/python-service
autorestart=true
startsecs=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=CORS_ORIGIN="%(ENV_CORS_ORIGIN)s",HOST="127.0.0.1",PORT="5001"

[program:backend]
command=node /app/backend/server.cjs
autorestart=true
startsecs=3
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PORT="3001",PYTHON_URL="http://127.0.0.1:5001",PYTHON_TIMEOUT_MS="%(ENV_PYTHON_TIMEOUT_MS)s",CORS_ORIGIN="%(ENV_CORS_ORIGIN)s",NODE_ENV="production"

[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
